# concentration.timer
[Unit]
Description=Timer for concentration.service

[Timer]
OnBootSec=1min
OnUnitActiveSec=10min

[Install]
WantedBy=concentration.service